# 8.1/devfs/etc/supervisord/conf.d/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid

loglevel=debug
user=root

[program:fpm]
command=/usr/sbin/php-fpm81 -R --nodaemonize

[program:ssh]
command=/usr/sbin/sshd -D
